# systemd/mirrorlist.path
[Path]
PathExists=/etc/pacman.d/mirrorlist.pacnew

[Install]
WantedBy=multi-user.target
